# docker/php-cli/supervisor.d/conf.d/messenger-async.conf
[program:messenger-async]
command=/var/www/bin/console messenger:consume async --time-limit=3600 -vv
directory=/var/www
autostart=true
autorestart=true
startretries=5
startsecs=1
user=root
numprocs=1
stopsignal=KILL
process_name=%(program_name)s_%(process_num)02d
stderr_logfile=/var/log/supervisor/%(program_name)s_stderr.log
stderr_logfile_maxbytes=10MB
stdout_logfile=/var/log/supervisor/%(program_name)s_stdout.log
stdout_logfile_maxbytes=10MB